# docker-cypress/config/supervisord-httpd.conf
[program:httpd]
command=httpd -D "FOREGROUND" -k start 
autostart=true
autorestart=true