# docker-cypress/config/supervisord-mysqld.conf
[program:mysqld]
command=/usr/sbin/mysqld -u root
autostart=true
autorestart=true