# docker-cypress/config/supervisord-php-fpm.conf
[program:php-fpm]
command=/usr/sbin/php-fpm -R -F -c /etc/php/fpm/php-fpm.conf
autostart=true
autorestart=true